# smcroute.service.in
[Unit]
Description=Static multicast routing daemon
Documentation=man:smcroute
Documentation=file:@DOCDIR@/README.md
Documentation=https://github.com/troglobit/smcroute
# ConditionPathExists=/etc/smcroute.conf
After=network-online.target
Requires=network-online.target

[Service]
Type=notify
ExecStart=@SBINDIR@/smcrouted -n -s
ExecReload=@SBINDIR@/smcroutectl restart

# Hardening settings
NoNewPrivileges=true
ProtectControlGroups=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
